# udev-rules/20-external-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/su goy -c /home/goy/.local/bin/monitor-control"
